# common/.config/profile.d/70-locale.sh
#!/bin/sh
# ~/.config/profile.d/70-locale.sh
# set locale if not defined

[ -z "${LANG}" ] && LANG=en_US.UTF-8
