# common/.config/profile.d/80-aliases.sh
#!/bin/sh
# ~/.config/profile.d/80-aliases.sh

# all aliases for shell
